# deploy/systemd/anki-sync.service
[Unit]
Description=Trigger AnkiWeb sync through AnkiConnect
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=%h/anki-desktop-docker
ExecStart=/usr/bin/env bash %h/anki-desktop-docker/scripts/sync-now.sh
