# deploy/systemd/anki-sync.timer
[Unit]
Description=Run Anki sync every 10 minutes

[Timer]
OnBootSec=2m
OnUnitActiveSec=10m
Persistent=true

[Install]
WantedBy=timers.target
